# etc/systemd/user/tailscale-online.target
[Unit]
Description=Synced with Tailscale status