# etc/systemd/system/tailscale-online.target
[Unit]
Description=Synced with Tailscale status
BindsTo=tailscaled.service
After=tailscaled.service