# etc/systemd/system/tailscale-online@.target
[Unit]
Description=Synced with Tailscale status of the %i node
BindsTo=tailscaled.service
After=tailscaled.service